# gateway-agent/systemd/plg-relay.service
[Unit]
Description=PLGames PLG Relay Server
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/plg-relay
Restart=always
RestartSec=3
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
